# base/s6-rc.d/oci-env-prepare/up
#!/usr/bin/execlineb -P

bash /utils.sh init_container
