# base/s6-rc.d/oci-env-prepare/up
#!/usr/bin/execlineb -P
foreground { s6-echo "Calling /utils.sh init_container" }

bash /utils.sh init_container
