# base/s6-rc.d/profile-init/up
#!/usr/bin/execlineb -P

with-contenv

execline-cd /
bash /src/oci_env/.compiled/init.sh
