# base/s6-rc.d/oci-env-profiles/up
#!/usr/bin/execlineb -P

bash /src/oci_env/.compiled/init.sh