# base/s6-rc.d/oci-env-profiles/up
#!/usr/bin/execlineb -P

foreground { s6-echo "Calling /utils.sh run_profile_init_scripts" }

bash /utils.sh run_profile_init_scripts