# base/s6-rc.d/oci-env-profiles/up
#!/usr/bin/execlineb -P

foreground { s6-echo "Calling /src/oci_env/.compiled/init.sh" }

bash /src/oci_env/.compiled/init.sh